# files/userconfig/autostart/albert.desktop
[Desktop Entry]
Name=Synapse
Exec=albert
Encoding=UTF-8
Type=Application
Icon=albert
Comment=
RunHook=0

