# files/userconfig/autostart/synapse.desktop
[Desktop Entry]
Name=Synapse
Exec=sh -c "sleep 10; synapse --startup"
Encoding=UTF-8
Type=Application
Icon=synapse
Comment=
RunHook=0

